# magento/crontab
*/1 * * * * php -c /etc/ /var/www/magento2/bin/magento cron:run
*/1 * * * * php -c /etc/ /var/www/magento2/update/cron.php
*/1 * * * * php -c /etc/ /var/www/magento2/bin/magento setup:cron:run
